# alert-notification@.service
[Unit]
Description=triggering alert for %i

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /opt/systemd-alerts/systemd-alerts.sh %i
