# roles/systemd-alerts/templates/alert-notification@.service.j2
[Unit]
Description=triggering alert for %i
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash {{ system.install_directory }}/systemd-alerts.sh %i
